# config/s6/services/webman/run
#!/bin/sh
# Webman service - S6 Overlay

cd /app
exec php start.php start
